# nginx/root/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash

ls -lah /usr/sbin

exec /usr/sbin/nginx