# nginx/root/etc/services.d/dnsmasq/run
#!/usr/bin/with-contenv bash

exec dnsmasq --listen-address=127.0.0.1 -d